# deploy_tools/gunicorn-systemd.template.service
[Unit]
Description=Gunicorn server for SITENAME

[Service]
Restart=on-failure
User=liyanjie
WorkingDirectory=/home/liyanjie/sites/SITENAME/source
ExecStart=/home/liyanjie/sites/SITENAME/virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket superlists.wsgi:application

[Install]
WantedBy=multi-user.target
